# archive/services/launchagents/com.ethan.longrein-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ethan.longrein-sync</string>
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/node</string>
        <string>/opt/homebrew/bin/longrein</string>
        <string>install</string>
        <string>--link</string>
        <string>-y</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>StartInterval</key>
    <integer>600</integer>
    <key>WatchPaths</key>
    <array>
        <string>/Users/ethan/skills-hubs/longrein/skills</string>
        <string>/Users/ethan/skills-hubs/longrein/global</string>
    </array>
    <key>StandardOutPath</key>
    <string>/Users/ethan/.longrein/sync.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/ethan/.longrein/sync.err.log</string>
</dict>
</plist>
